# edge/service-templates/linux/edge-sync.timer
[Unit]
Description=GS Edge Sync Timer

[Timer]
OnBootSec=60
OnUnitActiveSec=300

[Install]
WantedBy=timers.target
